# crontab
*/5 * * * * ~/websites/CV/DeployUpdate.sh > ~/logs/crontab/CV_Deploy.txt